# crontab
* * * * * curl https://oyster-app-atyts.ondigitalocean.app/api/schedule/1X91H0c3Cg >/proc/1/fd/1 2>/proc/1/fd/2



# An empty line is required at the end of this file for a valid cron file.
